# systemd/ecclesiastic-prefetch.service
[Unit]
Description=Prefetch art wallpapers from Met Museum API
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/ecclesiastic-prefetch
Environment=HOME=%h
TimeoutStartSec=600

[Install]
WantedBy=default.target
